# crontab
@reboot cd /home/pi/lepton-mobile-camera && python camera.py &
